# assets/getty/etc/systemd/system/serial-getty@ttyS1.service.d/override.conf
# Provide a root shell on the secondary serial device
# (the one connected to shell.sock on the outside).

[Service]
ExecStart=
ExecStart=-/sbin/agetty --autologin root --noclear %I $TERM
