# assets/gitlab-runner/etc/systemd/system/clear-docker-cache.service
[Unit]
Description=Clean up old docker images

[Service]
Type=oneshot
ExecStart=/usr/share/gitlab-runner/clear-docker-cache
