# assets/gitlab-runner/etc/systemd/system/clear-docker-cache.timer
[Unit]
Description=Periodically clean up old docker images

[Timer]
OnCalendar=Thu, 04:33

[Install]
WantedBy=timers.target
